# devices/xunlong-orangepi3-lts/files/usr/lib/systemd/system/systemd-repart.service.d/root-device.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/repart-root-disk
