# devices/fairphone-fp5/files/etc/systemd/system/iio-sensor-proxy.service.d/hexagonrpcd.conf
# Make iio-sensor-proxy depend on hexagonrpcd
[Unit]
Wants=hexagonrpcd-adsp-sensorspd.service
After=hexagonrpcd-adsp-sensorspd.service
