# devices/fairphone-fp5/files/usr/lib/systemd/system/iio-sensor-proxy.service.d/10-fairphone-fp5.conf
[Unit]
Wants=hexagonrpcd-adsp-sensorspd.service
After=hexagonrpcd-adsp-sensorspd.service
After=fairphone-fp5-sensors.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
ExecStartPre=/bin/sleep 5
Restart=on-failure
RestartSec=10
RestrictAddressFamilies=AF_UNIX AF_LOCAL AF_NETLINK AF_QIPCRTR
