# devices/fairphone-fp5/files/usr/lib/systemd/system/iio-sensor-proxy.service.d/10-fairphone-fp5.conf
# Fairphone 5 override for iio-sensor-proxy
# Ensures sensor initialization happens before proxy starts

[Unit]
After=fairphone-fp5-sensors.service
Wants=fairphone-fp5-sensors.service

[Service]
# Add a delay to ensure ADSP/SSC is fully initialized
ExecStartPre=/bin/sleep 3
# Increase restart delay if service fails
RestartSec=5
# Environment for debugging (remove in production)
Environment="G_MESSAGES_DEBUG=all"