# devices/fairphone-fp5/files/usr/lib/systemd/system/iio-sensor-proxy.service.d/10-fairphone-fp5.conf
[Unit]
Wants=hexagonrpcd-adsp-sensorspd.service
After=hexagonrpcd-adsp-sensorspd.service
After=systemd-tmpfiles-setup.service

[Service]
Restart=on-failure
RestartSec=5
