# devices/fairphone-fp5/files/usr/lib/systemd/system/iio-sensor-proxy.service.d/10-fairphone-fp5.conf
[Unit]
Wants=hexagonrpcd-adsp-rootpd.service
Wants=hexagonrpcd-adsp-sensorspd.service
Wants=hexagonrpcd-sdsp.service
After=hexagonrpcd-adsp-rootpd.service
After=hexagonrpcd-adsp-sensorspd.service
After=hexagonrpcd-sdsp.service
