# devices/fairphone-fp5/files/usr/lib/systemd/system/iio-sensor-proxy.service.d/90-fairphone-fp5.conf
[Unit]
After=hexagonrpcd-adsp-sensorspd.service
Wants=hexagonrpcd-adsp-sensorspd.service

[Service]
RestrictAddressFamilies=
