# devices/fairphone-fp5/files/usr/lib/systemd/system/iio-sensor-proxy.service.d/fairphone-fp5.conf
[Unit]
Wants=hexagonrpcd-adsp-sensorspd.service
After=hexagonrpcd-adsp-sensorspd.service
After=fairphone-fp5-sensors.service

[Service]
ExecStartPre=/bin/sleep 3
