# devices/xunlong-orangepi3-lts/files/usr/lib/systemd/system/grow-rootfs.service.d/after-repart.conf
[Unit]
After=systemd-repart.service
Requires=systemd-repart.service
