# devices/oneplus-sdm845/files/fwload/usr/lib/systemd/system/ModemManager.service.d/wait-for-rmtfs.conf
[Unit]
After=rmtfs.service
Requires=rmtfs.service
