# devices/fairphone-fp5/files/usr/lib/systemd/system/ModemManager.service.d/10-fairphone-fp5.conf
[Unit]
After=set-modem-sim-slot.service
Requires=set-modem-sim-slot.service
