# devices/fairphone-fp5/files/usr/lib/systemd/system/hexagonrpcd-adsp-sensorspd.service.d/fairphone-fp5.conf
[Unit]
After=sys-devices-platform-soc\x400-3700000.remoteproc.device
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Restart=on-failure
RestartSec=5
