# devices/fairphone-fp5/files/usr/lib/systemd/system/bluetooth.service.d/10-fairphone-fp5.conf
[Unit]
# Ensure bluetooth starts after the WCN6750 HCI device is available.
# The BT serdev driver downloads firmware during early boot and registers
# hci0. We want bluetoothd to start only after this is complete.
Wants=sys-subsystem-bluetooth-devices-hci0.device
After=sys-subsystem-bluetooth-devices-hci0.device

[Service]
# WCN6750 BT workaround: give the HCI device extra time after firmware
# download to complete post-setup initialization. The "unexpected event
# for opcode 0xfc48" seen during setup can delay readiness.
ExecStartPre=/bin/sleep 3
