# devices/fairphone-fp5/files/usr/lib/systemd/system/bluetooth.service.d/10-fairphone-fp5.conf
[Unit]
Requires=sys-subsystem-bluetooth-devices-hci0.device
After=sys-subsystem-bluetooth-devices-hci0.device

[Service]
# Match the directory mode that bluez package creates
ConfigurationDirectoryMode=0755
