# devices/fairphone-fp5/files/usr/lib/systemd/system/bluetooth.service.d/10-fairphone-fp5.conf
[Unit]
Requires=sys-subsystem-bluetooth-devices-hci0.device
After=sys-subsystem-bluetooth-devices-hci0.device

[Service]
ExecStartPre=/usr/libexec/fairphone-fp5-bt-setup
