# devices/fairphone-fp5/files/usr/lib/systemd/system/fairphone-fp5-bluetooth.service
[Unit]
Description=Fairphone 5 Bluetooth Initialization
Description=Initialize Bluetooth on Fairphone 5 (WCN6750 via WCNSS)
After=bluetooth.target
After=sys-subsystem-net-devices-wlan0.device
Wants=bluetooth.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 10
ExecStart=/usr/libexec/fairphone-fp5-bluetooth-init
TimeoutStartSec=180
Restart=on-failure
RestartSec=10

[Install]
WantedBy=bluetooth.target
WantedBy=multi-user.target
