# devices/fairphone-fp5/files/usr/lib/systemd/system/fairphone-fp5-bluetooth.service
[Unit]
Description=Fairphone 5 Bluetooth Initialization (WCN6750)
After=sys-subsystem-bluetooth-devices-hci0.device
Before=bluetooth.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/fairphone-fp5-bluetooth-init
ExecStartPost=/bin/sh -c 'systemctl try-restart bluetooth.service 2>/dev/null || true'
TimeoutStartSec=60

[Install]
WantedBy=bluetooth.target
