# devices/fairphone-fp5/files/usr/lib/systemd/system/set-bt-mac.service
[Unit]
Description=Set Bluetooth MAC address

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/set-bt-mac
