# devices/fairphone-fp5/files/usr/lib/systemd/system/ucsi-rebind.service
[Unit]
Description=Rebind UCSI to fix USB 2.0 HS with DP alt-mode partners

[Service]
Type=oneshot
ExecStart=/usr/libexec/ucsi-rebind
